# scripts/clean_downloads.service
[Unit]
Description=Clean Esme's Downloads folder

[Service]
Type=oneshot
ExecStart=%h/dotfiles/scripts/clean_downloads.sh